# production/etc/service/bevly-server/run
#!/bin/sh
exec /sbin/setuser nobody /var/lib/bevly/bevly-server