# production/etc/service/bevly-server/run
#!/bin/sh
export MONGO_HOST=mongodb
exec /sbin/setuser nobody /var/lib/bevly/bevly-server >/var/log/bevly.log 2>&1
