# ubuntu-init/bevly-server.conf
# bevly-server - run bevly beer server docker container

description "Bevly server (docker)"

start on started docker
task

script
  docker start bevly-docker
end script
