# production/etc/supervisor/conf.d/bevly-server.conf
[program:bevly]
command=/root/bevly/bevly-server
user=nobody