# production/etc/supervisor/conf.d/bevly-server.conf
[program:bevly]
command=/var/lib/bevly/bevly-server
user=nobody